# scripts/work-parser-run-all.service
[Unit]
Description=Run all parser strategies pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=oleg
WorkingDirectory=/home/oleg/parser/work_parser

# Optional environment variables file.
EnvironmentFile=-/home/oleg/parser/work_parser/.env

# Optional arguments passed from env file, for example:
# RUN_ALL_ARGS="-m 200 -w 10,110,210 --skip-generate"
Environment=RUN_ALL_ARGS=

# Main pipeline script (independent from current terminal session).
ExecStart=/bin/bash -lc '/home/oleg/parser/work_parser/scripts/run-all-strategies.sh ${RUN_ALL_ARGS}'

StandardOutput=journal
StandardError=journal
TimeoutStartSec=12h

[Install]
WantedBy=multi-user.target
