# scripts/work-parser-run-all@.service
[Unit]
Description=Run parser pipeline with args: %I
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=oleg
WorkingDirectory=/home/oleg/parser/work_parser
EnvironmentFile=-/home/oleg/parser/work_parser/.env

# Start instance with escaped arguments in %I, for example:
# systemctl start "work-parser-run-all@-m\x20100\x20--skip-generate.service"
ExecStart=/bin/bash -lc '/home/oleg/parser/work_parser/scripts/run-all-strategies.sh %I'

StandardOutput=journal
StandardError=journal
TimeoutStartSec=12h

[Install]
WantedBy=multi-user.target
